# salt/roots/nodejs/supervisor.conf
[program:LeGrogJdRA]
command=grunt
directory=/home/legrog/LeGrogJdRA/
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/LeGrogJdRA/LeGrogJdRA.err.log
stdout_logfile=/var/log/LeGrogJdRA/LeGrogJdRA.out.log
user=legrog
environment=NODE_ENV='production'